# scripts/openenvx-orchestrator.service
[Unit]
Description=OpenEnvX Orchestrator Daemon
After=network.target postgresql.service

[Service]
ExecStart=/usr/local/bin/orchestrator
Restart=always
EnvironmentFile=/etc/openenvx/orchestrator.env
TimeoutStopSec=30
User=openenvx

[Install]
WantedBy=multi-user.target
